# scripts/zapme.service
# Systemd user-level unit for zapme — vision-driven posture corrector.
#
# Boots the runtime loop automatically when the Pi powers on.
# Combined with `loginctl enable-linger` (set automatically by the
# install script), the user's session — and therefore this service —
# starts at boot before anyone has to log in.
#
# Why user-level instead of system-level:
#   - Audio (PipeWire / PulseAudio / ALSA fallback) runs inside the
#     user session. A system-level service has no audio context at
#     boot.
#   - Group membership (gpio / audio / video) is automatic — the
#     service runs as the user and inherits all their groups.
#   - No __USER__ template substitution: user services are owned by
#     whichever user installs them.
#
# Install (renders nothing — just copies into the user's systemd dir):
#
#   sudo bash scripts/install_systemd.sh
#
# That script also calls `loginctl enable-linger`, which is what
# makes the service start at boot without a login.
#
# Manual control after install:
#
#   systemctl --user status  zapme.service
#   systemctl --user restart zapme.service
#   journalctl --user -u zapme.service -f
#
# Stop / disable / uninstall:
#
#   systemctl --user stop    zapme.service
#   systemctl --user disable zapme.service
#   rm ~/.config/systemd/user/zapme.service
#   systemctl --user daemon-reload
#   sudo loginctl disable-linger $USER     # if you don't want the
#                                          # user session at boot anymore
#
# Edit GPIO pin, camera index, thresholds etc. directly in the
# ExecStart= block below if they differ from the defaults.

[Unit]
Description=Zapme posture corrector
# default.target is the user-session equivalent of multi-user.target.
After=default.target
# Explicit dependency on the audio stack so we don't race PipeWire /
# WirePlumber at boot. Without these, on Pi OS Trixie the service
# can start before the audio socket exists and silently no-op all
# the .play() calls. install_systemd.sh enables these user units
# alongside zapme so they actually exist to depend on.
Wants=pipewire.service pipewire-pulse.service wireplumber.service
After=pipewire.service pipewire-pulse.service wireplumber.service

[Service]
Type=simple
WorkingDirectory=/opt/zapme

# Force unbuffered stdout/stderr so journalctl shows logs in real time
Environment=PYTHONUNBUFFERED=1

# Use the venv's Python directly. The venv must be created with
# `--system-site-packages` so it can see the apt-installed lgpio.
ExecStart=/opt/zapme/.venv/bin/python -m zapme.src \
    --backend lgpio \
    --weights /opt/zapme/models/slouch_cnn.onnx \
    --camera 0 \
    --gate-pin 17 \
    --cooldown 15 \
    --pulse-duration 0.3 \
    --on-threshold 0.8 \
    --off-threshold 0.4
# Polarity defaults to active-LOW (line LOW = relay ON, HIGH at boot
# = de-energized). If you ever swap to a relay that's wired the other
# way (active-HIGH), append `--gate-active-high`. Verify with:
#   /opt/zapme/.venv/bin/python -m zapme.src.utils.tens_control --pulse 0.5
# A clean silent→click→silent (with EMS pads unplugged) confirms the
# default is correct.

# Restart only on non-zero exit (watchdog trip → exit 1). Don't restart
# on clean exit (operator pressed Ctrl-C, exit 0).
Restart=on-failure
RestartSec=3

# Cap restart attempts so a stuck failure mode doesn't hammer the gate
# every few seconds. Five failures within a minute → stop trying.
StartLimitBurst=5
StartLimitIntervalSec=60

# Send SIGTERM on stop so the Python signal handler runs (closes the
# camera, stops the watchdog, drives the gate low). Wait up to 5s for
# clean shutdown before SIGKILL.
KillSignal=SIGTERM
TimeoutStopSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
