# scripts/zapme.service
# Systemd unit for zapme — vision-driven posture corrector.
#
# Boots the runtime loop automatically on Pi power-on. Combined with a
# physical power switch on the box, this gives the operator a single
# physical gesture to start / stop the whole system: flip the switch
# on, the Pi boots and zapme starts; flip it off, everything dies.
#
# Repo install location: /opt/zapme (fixed). Only `User=` is per-machine
# and gets substituted at install time. See `scripts/install_systemd.sh`
# for the full one-shot install, or run the steps manually:
#
#   sudo mkdir -p /opt/zapme && sudo chown -R $USER:$USER /opt/zapme
#   git clone <repo-url> /opt/zapme       # or `mv ~/zapme /opt/zapme`
#   cd /opt/zapme
#   sudo apt install -y mpg123 python3-lgpio
#   python3 -m venv .venv --system-site-packages
#   .venv/bin/pip install -e .[pi]
#
#   sed "s|__USER__|$USER|g" scripts/zapme.service \
#     | sudo tee /etc/systemd/system/zapme.service > /dev/null
#   sudo systemctl daemon-reload
#   sudo systemctl enable  zapme.service      # auto-start on boot
#   sudo systemctl start   zapme.service      # start now
#   sudo systemctl status  zapme.service      # check
#   journalctl -u zapme.service -f            # live logs
#
# Stop / disable / uninstall:
#
#   sudo systemctl stop    zapme.service
#   sudo systemctl disable zapme.service
#   sudo rm /etc/systemd/system/zapme.service
#   sudo systemctl daemon-reload
#
# Edit GPIO pin, camera index, thresholds etc. directly in the
# ExecStart= block below if they differ from the defaults.

[Unit]
Description=Zapme posture corrector
# Wait for udev so /dev/video* and /dev/gpiochip* are populated
After=systemd-udev-settle.service
Wants=systemd-udev-settle.service

[Service]
Type=simple
User=__USER__
Group=gpio
WorkingDirectory=/opt/zapme

# Force unbuffered stdout/stderr so journalctl shows logs in real time
Environment=PYTHONUNBUFFERED=1

# Audio routing for the system-level service. Pi OS Bookworm uses
# PipeWire / PulseAudio, which lives in the user's runtime directory
# at /run/user/<UID>. %U expands to the User=__USER__'s UID.
# REQUIRES `sudo loginctl enable-linger __USER__` to be run once so
# the user session (and therefore /run/user/%U) exists at boot before
# anyone logs in.
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=PULSE_RUNTIME_PATH=/run/user/%U/pulse

# Use the venv's Python so the service does not depend on system Python
# being correctly configured. The venv must be created with
# `--system-site-packages` so it can see the apt-installed lgpio.
ExecStart=/opt/zapme/.venv/bin/python -m zapme.src \
    --backend lgpio \
    --weights /opt/zapme/models/slouch_cnn.onnx \
    --camera 0 \
    --gate-pin 17 \
    --cooldown 15 \
    --pulse-duration 0.3 \
    --on-threshold 0.8 \
    --off-threshold 0.4
# Polarity defaults to active-LOW (line LOW = relay ON, HIGH at boot
# = de-energized). If you ever swap to a relay that's wired the other
# way (active-HIGH), append `--gate-active-high`. Verify with:
#   /opt/zapme/.venv/bin/python -m zapme.src.utils.tens_control --pulse 0.5
# A clean silent→click→silent (with EMS pads unplugged) confirms the
# default is correct.

# Restart only on non-zero exit (watchdog trip → exit 1). Don't restart
# on clean exit (operator pressed Ctrl-C, exit 0).
Restart=on-failure
RestartSec=3

# Cap restart attempts so a stuck failure mode doesn't hammer the gate
# every few seconds. Five failures within a minute → stop trying.
StartLimitBurst=5
StartLimitIntervalSec=60

# Send SIGTERM on stop so the Python signal handler runs (closes the
# camera, stops the watchdog, drives the gate low). Wait up to 5s for
# clean shutdown before SIGKILL.
KillSignal=SIGTERM
TimeoutStopSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
